# launchagents/com.example.relay-drones.triage.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.relay-drones.triage</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{PY}}</string>
        <string>{{REPO}}/bin/relay-drones</string>
        <string>triage</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{REPO}}</string>

    <key>RunAtLoad</key><false/>
    <key>StartInterval</key><integer>300</integer>

    <key>StandardOutPath</key>
    <string>{{REPO}}/logs/triage.log</string>
    <key>StandardErrorPath</key>
    <string>{{REPO}}/logs/triage.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key><string>1</string>
    </dict>
</dict>
</plist>
